# config/sys/fwcloud-websrv.service
[Unit]
Description=fwcloud-websrv
After=fwcloud-api.service

[Service]
ExecStart=/usr/bin/npm start
ExecStop=/usr/bin/npm run stop
WorkingDirectory=/opt/fwcloud/websrv/
SyslogIdentifier=fwcloud-websrv
User=fwcloud
Group=fwcloud
Type=simple

[Install]
WantedBy=multi-user.target
Alias=fwcloud-ui.service